# bgrun.service
[Unit]
Description=bgrun daemon

[Service]
ExecStart=/usr/local/bin/bgrun --daemon
Restart=always

[Install]
WantedBy=default.target
