# bgrun.service
[Unit]
Description=bgrun daemon

[Service]
ExecStart=/usr/local/bin/bgrun -s
Restart=always

[Install]
WantedBy=default.target
